# ops/systemd/company-agent-sessions-exporter.timer
[Unit]
Description=Company agent-sessions exporter timer — every 20 min (exporter quiesce window is 15 min; cadence per Session Fabric synthesis Round 1b)

[Timer]
OnCalendar=*:00/20
Persistent=true
Unit=company-agent-sessions-exporter.service

[Install]
# Bound to the Company target (not bare timers.target) so it rises/falls with the
# Company and `company status` shows it as a managed member — Tim's requirement:
# no invisible orphan timer. The TIMER (not the oneshot .service) is the member,
# matching how ops/services.json declares this job's managed unit.
WantedBy=company.target
