# ops/systemd/company-config-writer.service
[Unit]
Description=Company config writer — R3 .claude/config/VCS mutator (HTTP on 127.0.0.1:8772)
After=network.target
PartOf=company.target

[Service]
Type=simple
WorkingDirectory=/home/tim/company
ExecStart=/home/tim/company/.venv/bin/python runtime/config_writer.py 8772
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=company-config-writer

[Install]
WantedBy=company.target
