# ops/systemd/litellm-proxy.service
[Unit]
Description=LiteLLM Proxy Gateway for Claude Code
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/claude-gateway
ExecStart=%h/claude-gateway/.venv/bin/litellm --config %h/claude-gateway/config.yaml --port 4000
EnvironmentFile=%h/claude-gateway/.env
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
