# ops/systemd/vllm-chat-4b-fp8.service
[Unit]
Description=vLLM Chat Server (Qwen3.5-4B-FP8, RedHatAI) on port 8001
After=default.target

[Service]
Type=simple
ExecStart=/home/tim/company/ops/serve_model.sh chat-4b-fp8
Restart=on-failure
RestartSec=10
# Logs go to journalctl: journalctl --user -u vllm-chat-4b-fp8 -f
StandardOutput=journal
StandardError=journal
# Generous timeout for cold-start compilation
TimeoutStartSec=300

[Install]
WantedBy=default.target
