# ops/systemd/vllm-chat.service
[Unit]
Description=vLLM Chat Server (Qwen3.5-4B-AWQ) on port 8000
After=default.target

[Service]
Type=simple
ExecStart=/home/tim/company/ops/serve_model.sh chat-4b
Restart=on-failure
RestartSec=10
# Logs go to journalctl: journalctl --user -u vllm-chat -f
StandardOutput=journal
StandardError=journal
# Generous timeout for cold-start compilation
TimeoutStartSec=300

[Install]
WantedBy=default.target
