# voice/ops/systemd/company-voice-chatterbox.service
[Unit]
Description=Company voice engine — Chatterbox-TTS on port 4124
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/tim/company
EnvironmentFile=%h/company/voice/ops/voice.env
ExecStart=%h/.voice-venvs/chatterbox/bin/python /home/tim/company/voice/engines/chatterbox.py 4124
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
# Generous: first run downloads the HF model + loads onto GPU.
TimeoutStartSec=600

[Install]
WantedBy=default.target
